# salt/roles/saltmaster/files/saltmaster.cron
*/5 * * * * root chown -R root:salt /srv/salt
*/5 * * * * root find /srv/salt -type d -exec chmod 2770 -- {} +
*/5 * * * * root find /srv/salt -type f -exec chmod 0660 -- {} +
