# salt/roles/vault/files/vault.conf
{% from 'roles/vault/map.jinja' import vault with context -%}

[program:vault]
command={{ vault.base_dir }}/bin/vault-{{ vault.version }} server -config {{ vault.base_dir }}/config -log-level=info
stdout_logfile={{ vault.base_dir }}/log/vault-stdout.log
stderr_logfile={{ vault.base_dir }}/log/vault-stderr.log
numprocs=1
autostart=true
autorestart=true
stopsignal=INT
user={{ vault.user }}
